# scripts/automation/com.dotfiles.serena-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dotfiles.serena-update</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>%HOME%/dotfiles/scripts/serena-auto-update.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>%HOME%/dotfiles</string>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>StandardOutPath</key>
    <string>%HOME%/dotfiles/.serena/automation.log</string>
    
    <key>StandardErrorPath</key>
    <string>%HOME%/dotfiles/.serena/automation_error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>